# .claude/worktrees/carolina-node24-npm-audit/systemd/isnad-backup.service
[Unit]
Description=isnad-graph automated database backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/isnad-graph/scripts/backup.sh
WorkingDirectory=/opt/isnad-graph
EnvironmentFile=/opt/isnad-graph/.env
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/tmp/isnad-backups /opt/isnad-graph
PrivateTmp=yes
