# .claude/worktrees/carolina-node24-npm-audit/systemd/isnad-backup.timer
[Unit]
Description=Daily backup timer for isnad-graph databases

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
